# gunicorn.conf
[program:myserver]
command=gunicorn knowgit:app -w 4 -b :80

[supervisord]
logfile=/home/ubuntu/supervisord.log
loglevel=debug
user=root
